# .config/systemd/user/i3lock.service
[Unit]
Description=Lock screen before sleep

[Service]
Type=oneshot
ExecStar=/usr/bin/i3lock
RemainAfterExit=yes

[Install]
WantedBy=sleep.target
